# server-prod/meter-audit.timer
[Unit]
Description=Run the water meter accuracy audit every 20 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=20min
Persistent=true

[Install]
WantedBy=timers.target
