# ops/codex-meter-improvement/smart-garden-meter-improvement.timer
[Unit]
Description=Run meter OCR engineering agent daily

[Timer]
OnCalendar=*-*-* 03:38:50
Persistent=true
AccuracySec=10s
Unit=smart-garden-meter-improvement.service

[Install]
WantedBy=timers.target
