# cnn/meter-cnn-retrain.service
[Unit]
Description=Nightly gated retrain of the meter-digit CNN (champion/challenger)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=jack
WorkingDirectory=/home/jack/meter-cnn
# Train politely so the live CNN service + RapidOCR stay responsive.
# This unit may run frequently; retrain.py exits quickly when MIN_NEW_FRAMES
# is not met, so frequent checks stay cheap.
Nice=10
IOSchedulingClass=idle
ExecStart=/home/jack/meter-ocr/.venv/bin/python /home/jack/meter-cnn/retrain.py
TimeoutStartSec=14400
