# cnn/meter-cnn-retrain.timer
[Unit]
Description=Run the meter-CNN gated retrain frequently

[Timer]
# Near-real-time retrain checks: run often; retrain.py still gates on
# MIN_NEW_FRAMES, so most cycles are lightweight skips until enough new data
# arrives.
OnUnitActiveSec=10min
# If the box was asleep/off, run soon after boot instead of waiting for a
# clock time.
OnBootSec=2min
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
